# stemcell_builder/stages/system_oneandone_packages/assets/privatenet.conf
#Start the private netowrk configuration service


start on startup and net-device-up IFACE=eth0

task
script
    exec  /usr/share/scripts/privatenet.sh
end script